# config/supervisord.conf
[supervisord]
nodaemon=true

[program:kodi]
command=LD_PRELOAD="/opt/vc/lib/libGLESv2.so /opt/vc/lib/libEGL.so" /usr/bin/kodi-standalone
